# bashrc/available-upgrades/available-upgrades.timer
[Unit]
Description=Periodically check for available package upgrades

[Timer]
OnBootSec=15min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
